# .config/autostart/xbindkeys.desktop
[Desktop Entry]
Type=Application
Name=xbindkeys
TryExec=xbindkeys
Exec=xbindkeys -n
StartupNotify=false
